# deploy/systemd/bilson-ai.service
# Bilson AI web service (website + accounts + API keys + admin + REST API).
# Loads the search model once and serves everything on 127.0.0.1:8001; Caddy
# terminates TLS in front. Install:
#   sudo cp deploy/systemd/bilson-ai.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now bilson-ai
# Edit User=, the BILSON_SECRET, BILSON_ADMIN_EMAIL, and paths below.
[Unit]
Description=Bilson AI web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anglican
WorkingDirectory=/home/anglican/anglican-search-starter
Environment=BILSON_HOST=127.0.0.1
Environment=BILSON_PORT=8001
# REQUIRED: set to the output of `openssl rand -hex 32` (keeps sessions valid across restarts).
Environment=BILSON_SECRET=CHANGE_ME
Environment=BILSON_DB=/home/anglican/anglican-search-starter/data/accounts.db
Environment=BILSON_MONTHLY_LIMIT=1000
# This email becomes an admin when it signs up (or use the admin CLI).
Environment=BILSON_ADMIN_EMAIL=you@example.com
# Public base URL (so the dashboard shows the correct MCP endpoint, e.g. .../mcp).
Environment=BILSON_PUBLIC_URL=https://your.domain.example
# Search engine (shared, loaded once). On the RTX 3060 box we use the full
# GPU reranker. (For a CPU-only host instead: ANGLICAN_RERANKER=cross-encoder/
# ms-marco-MiniLM-L-6-v2 and ANGLICAN_RERANK_POOL=30.)
Environment=ANGLICAN_DB=/home/anglican/anglican-search-starter/data/library.db
Environment=ANGLICAN_INDEX=/home/anglican/anglican-search-starter/data/index.faiss
Environment=ANGLICAN_RERANKER=BAAI/bge-reranker-v2-m3
Environment=ANGLICAN_RERANK_POOL=80
ExecStart=/home/anglican/.local/bin/uv run bilson-ai
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
